# systemd/codeagent.service
[Unit]
Description=CodeAgent-RWKV API Server
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/codeagent-server
WorkingDirectory=%h/.local/share/codeagent
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
Environment=HOST=127.0.0.1
Environment=PORT=8080
Environment=MAX_TOKENS=512
Environment=TEMPERATURE=0.7
Environment=TOP_P=0.9
Environment=CHECKPOINT=%h/.local/share/codeagent/checkpoints
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
